# deploy/systemd/llm-search-mediator.service
[Unit]
Description=LLM Search Mediator (privacy-preserving search bridge)
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m search_mediator.app
WorkingDirectory=/opt/llm-search-mediator
Environment=BIND_ADDR=127.0.0.1:8485
Environment=SEARXNG_URL=http://127.0.0.1:8888
Environment=POLICY_PATH=/etc/llm-search-mediator/policy.yaml
Environment=AUDIT_DIR=/var/lib/llm-search-mediator/logs

# Filesystem isolation
DynamicUser=yes
StateDirectory=llm-search-mediator
LogsDirectory=llm-search-mediator
ReadWritePaths=/var/lib/llm-search-mediator
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ProtectProc=invisible
ProcSubset=pid

# Kernel protection
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes

# Privilege restriction
NoNewPrivileges=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes
RestrictRealtime=yes

# Capability bounding
CapabilityBoundingSet=
AmbientCapabilities=

# Namespace restrictions
RestrictNamespaces=yes
RestrictAddressFamilies=AF_INET AF_UNIX
PrivateUsers=yes

# Syscall filtering
SystemCallFilter=@system-service @network-io
SystemCallFilter=~@privileged @resources @mount @clock @debug @swap @reboot @raw-io @module @cpu-emulation @obsolete
SystemCallArchitectures=native
SystemCallErrorNumber=EPERM

# Resource limits
MemoryMax=256M
CPUQuota=25%
TasksMax=64
LimitNOFILE=512
LimitCORE=0

Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
